# systemd/.config/systemd/user/mail-sync.timer
[Unit]
Description=Gets the mail

[Timer]
OnStartupSec=2m
OnUnitActiveSec=5m
WakeSystem=true

[Install]
WantedBy=timers.target
